# scripts/systemd/coach-web.service
[Unit]
Description=Coach Tools — Next.js web app
After=network.target

[Service]
Type=simple
User=cqlanus
WorkingDirectory=/home/cqlanus/Documents/coach-tools/apps/web
ExecStart=/home/cqlanus/.nvm/versions/node/v24.16.0/bin/node node_modules/.bin/next start -p 3010
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
Environment=NEXT_PUBLIC_UMAMI_WEBSITE_ID=793e2f93-bcb1-484e-bea5-918d77945be0

[Install]
WantedBy=multi-user.target
