# scripts/systemd/coach-api.service
[Unit]
Description=Coach Tools — FastAPI server
After=network.target

[Service]
Type=simple
User=cqlanus
WorkingDirectory=/home/cqlanus/Documents/coach-tools/api
ExecStart=/home/cqlanus/Documents/coach-tools/api/.venv/bin/uvicorn main:app --port 3011 --host 127.0.0.1
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
